# packaging/systemd/showeq-daemon.service
[Unit]
Description=ShowEQ packet-capture daemon
Documentation=https://github.com/showeq-unofficial/showeq-daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# All daemon flags come from this file so the unit itself doesn't need
# editing per host. See showeq-daemon.env.example for the keys.
EnvironmentFile=/etc/showeq-daemon/showeq-daemon.env

# Read-only opcode XML / spells_us.txt / eqstr_us.txt are picked up from
# the compiled-in PKGDATADIR (default /usr/local/share/showeq-daemon).
# Writable user state lands in $HOME/.showeq — the same dir legacy
# showeq writes, so filters/, maps/, spawnpoints/ etc. interop directly.
# Set User= below to your account so $HOME resolves to /home/<you>.
ExecStart=/usr/local/bin/showeq-daemon \
    --device=${SEQ_DEVICE} \
    --listen=${SEQ_LISTEN} \
    $SEQ_EXTRA_ARGS

# Run as your unprivileged account so $HOME = /home/<you> and
# DataLocationMgr writes to /home/<you>/.showeq (shared with legacy
# showeq). User=root is the easy-but-wrong default — under systemd
# there is no SUDO_USER, so root means /root/.showeq, no interop.
# pcap needs CAP_NET_RAW + CAP_NET_ADMIN; the ambient capabilities
# grant them without making the binary setuid root.
User=root
# Recommended:
# User=YOUR_ACCOUNT
# AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
# CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Crash-loop-friendly restart. EQ session crypto means a fresh daemon
# process can't decrypt mid-session traffic until the next login event
# (see NICE_TO_HAVE.md), so a brief gap on restart is expected — don't
# set RestartSec too low or systemd will spin if pcap is misconfigured.
Restart=on-failure
RestartSec=5

# stdout/stderr → journal by default; the in-process Qt logger already
# prepends ISO timestamps + level tags, which journalctl preserves.
StandardOutput=journal
StandardError=journal

# Modest hardening. ProtectHome is intentionally NOT set — the daemon
# needs to write into $HOME/.showeq for legacy-showeq interop.
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true

[Install]
WantedBy=multi-user.target
